# tcmu-runner.service
[Unit]
Description=LIO Userspace-passthrough daemon

[Service]
Type=dbus
BusName=org.kernel.TCMUService1
ExecStart=/usr/bin/tcmu-runner
